# npimp.service
[Unit]
Description=NTB pimp

[Service]
ExecStart=/usr/bin/npimp

[Install]
WantedBy=multi-user.target
